# contrib/dhcp-forwarder.service.in
# --*- conf -*--
[Unit]
Description = Agent to relay DHCP messages between two networks
After=syslog.target network.target

[Service]
Type = notify
ExecStart = @sbindir@/dhcp-fwd -c @sysconfdir@/dhcp-fwd.conf
Restart = on-failure

[Install]
WantedBy = multi-user.target
